# rootfs/etc/cron.d/certmaintain
00 * * * * root /bin/sh /usr/local/bin/certmaintain.sh > /dev/null 2>&1